# Tools/pkg/postinstall
#!/bin/bash
# NRIME PKG postinstall script

INSTALL_DIR="/Library/Input Methods"
APP_PATH="$INSTALL_DIR/NRIME.app"
SETTINGS_PATH="$INSTALL_DIR/NRIMESettings.app"
RESTORE_HELPER_PATH="$INSTALL_DIR/NRIMERestoreHelper.app"
LAUNCH_AGENT_PATH="/Library/LaunchAgents/com.nrime.inputmethod.loginrestore.plist"
MOZC_LAUNCH_AGENT_PATH="/Library/LaunchAgents/com.nrime.inputmethod.mozcserver.plist"
LSREGISTER="/System/Library/Frameworks/CoreServices.framework/Versions/Current/Frameworks/LaunchServices.framework/Versions/Current/Support/lsregister"

# Kill existing processes
killall NRIME 2>/dev/null || true
killall NRIMESettings 2>/dev/null || true
killall NRIMERestoreHelper 2>/dev/null || true
killall mozc_server 2>/dev/null || true
sleep 1

# Clean stale mozc lock files (read-only, chmod first)
CONSOLE_USER="$(stat -f "%Su" /dev/console 2>/dev/null || true)"
if [ -n "$CONSOLE_USER" ] && [ "$CONSOLE_USER" != "root" ]; then
    MOZC_DIR="$(eval echo "~$CONSOLE_USER")/Library/Application Support/Mozc"
    mkdir -p "$MOZC_DIR"
    chown "$CONSOLE_USER":staff "$MOZC_DIR" 2>/dev/null || true
    chmod 644 "$MOZC_DIR/.server.lock" 2>/dev/null || true
    rm -f "$MOZC_DIR/.server.lock" "$MOZC_DIR/.session.ipc" 2>/dev/null || true
fi

# Remove quarantine
xattr -cr "$APP_PATH" 2>/dev/null || true
xattr -cr "$SETTINGS_PATH" 2>/dev/null || true
xattr -cr "$RESTORE_HELPER_PATH" 2>/dev/null || true

# Register with LaunchServices
"$LSREGISTER" "$APP_PATH" 2>/dev/null || true
"$LSREGISTER" "$SETTINGS_PATH" 2>/dev/null || true
"$LSREGISTER" "$RESTORE_HELPER_PATH" 2>/dev/null || true

# Re-enable input sources (macOS auto-registers them from the bundle)
swift -e '
import Carbon
let conditions = [kTISPropertyBundleID: "com.nrime.inputmethod.app" as CFString] as CFDictionary
if let sources = TISCreateInputSourceList(conditions, true)?.takeRetainedValue() as? [TISInputSource] {
    for source in sources { TISEnableInputSource(source) }
}
' 2>&1 || true

# Reload login restore agent
if [ -n "$CONSOLE_USER" ] && [ "$CONSOLE_USER" != "root" ] && [ -f "$LAUNCH_AGENT_PATH" ]; then
    CONSOLE_UID="$(id -u "$CONSOLE_USER" 2>/dev/null || true)"
    if [ -n "$CONSOLE_UID" ]; then
        launchctl bootout "gui/$CONSOLE_UID" "$LAUNCH_AGENT_PATH" 2>/dev/null || true
        launchctl bootstrap "gui/$CONSOLE_UID" "$LAUNCH_AGENT_PATH" 2>/dev/null || true
    fi
fi

if [ -n "$CONSOLE_USER" ] && [ "$CONSOLE_USER" != "root" ] && [ -f "$MOZC_LAUNCH_AGENT_PATH" ]; then
    CONSOLE_UID="$(id -u "$CONSOLE_USER" 2>/dev/null || true)"
    if [ -n "$CONSOLE_UID" ]; then
        launchctl bootout "gui/$CONSOLE_UID" "$MOZC_LAUNCH_AGENT_PATH" 2>/dev/null || true
        launchctl bootstrap "gui/$CONSOLE_UID" "$MOZC_LAUNCH_AGENT_PATH" 2>/dev/null || true
        launchctl kickstart -k "gui/$CONSOLE_UID/com.nrime.inputmethod.mozcserver" 2>/dev/null || true
    fi
fi

exit 0
